# deploy/spring-backend.service
# Install: sudo cp deploy/spring-backend.service /etc/systemd/system/ && sudo systemctl daemon-reload && sudo systemctl enable --now spring-backend
# Edit paths if DEPLOY_PATH is not /home/srf. EnvironmentFile: same KEY=value as .env (no "export").

[Unit]
Description=Sword Of Knowledge Spring backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/srf
EnvironmentFile=-/home/srf/.env
# Match “ls /usr/lib/jvm” on this machine (amd64 vs arm64).
Environment=JAVA_HOME=/usr/lib/jvm/java-17-openjdk-amd64
ExecStart=/usr/bin/java -jar /home/srf/target/spring-backend-0.0.1-SNAPSHOT.jar
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
